# cron/crontab
* * * * * node /app/dist/cron/src/scheduledTasks/generateReport.js